# etc/systemd/system/cron.service.d/override.conf
[Service]
ExecStartPre=/bin/sh -c '/sbin/ip -o -4 addr show | grep -q "1.2.3.4" || (echo "Cron disable by /etc/systemd/system/cron.service.d/override.conf because of wrong IP"; /bin/false)'
